# roles/fedora/files/ssd/fstrim.cron
#!/bin/sh

#
# Trim free cells on the SSD
#
# Generated by Ansible
#
/sbin/fstrim --all
