# roles/mediacenter/templates/backup.service.j2
[Unit]
Description=Backup data for {{ item }}

[Service]
Type=oneshot
ExecStart=/usr/local/bin/{{ item }}-backup.sh
