# roles/mediacenter/templates/backup.timer.j2
[Unit]
Description=Backup data for {{ item }}

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
